# deploy/voicechanger.service
[Unit]
Description=Voice Changer real-time audio effect service
After=sound.target network.target
Documentation=https://github.com/example/voicechanger/wiki

[Service]
Type=simple
ExecStart=voicechanger serve --config ${HOME}/.config/voicechanger/voicechanger.toml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voicechanger
# Uncomment to run as dedicated user (requires user setup)
# User=voicechanger
# Group=audio
# Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=multi-user.target
